# ops/botty.service
[Unit]
Description=Botty application service
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=ofirkat
Group=ofirkat
SupplementaryGroups=docker
WorkingDirectory=/home/ofirkat/Botty
Environment=NODE_ENV=production
EnvironmentFile=/home/ofirkat/Botty/.env.local
ExecStartPre=/usr/bin/docker compose -f /home/ofirkat/Botty/docker-compose.yml up -d --wait ollama
ExecStart=/usr/bin/docker compose -f /home/ofirkat/Botty/docker-compose.yml up app
ExecStop=/usr/bin/docker compose -f /home/ofirkat/Botty/docker-compose.yml stop app
ExecStopPost=-/usr/bin/docker compose -f /home/ofirkat/Botty/docker-compose.yml rm -f app
Restart=always
RestartSec=5
TimeoutStartSec=0
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target