# ops/polymarket-adapter.service
[Unit]
Description=Polymarket Alpha Botty adapter (port 7002)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ofirkat
Group=ofirkat
WorkingDirectory=/home/ofirkat/eitan/polymarket-alpha
Environment=PATH=/home/ofirkat/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 /home/ofirkat/eitan/polymarket-alpha/botty_adapter.py
Restart=always
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
